# cron.d/tinymon
*/2 *     * * *     root   /usr/sbin/tinymon.sh
